# 3-discogs-site/x-nginx2/discogs-yt-playlister.conf
[program:discogs-yt-playlister]
command = gunicorn app:app -b localhost:8000
directory = /app
user = root